# imdb-relay.service
[Unit]
Description=IMDb List Relay API for Radarr
After=network.target

[Service]
User=matt
WorkingDirectory=/home/matt/gemini-projects/imdb_scraper_api
ExecStart=/home/matt/gemini-projects/imdb_scraper_api/.venv/bin/python3 main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
